# deploy/kalshi-mrb.service
[Unit]
Description=Kalshi Mean Reversion Bot
After=network.target

[Service]
Type=simple
User=andrew
WorkingDirectory=/opt/kalshi-mrb/backend
EnvironmentFile=/opt/kalshi-mrb/backend/.env
ExecStart=/opt/kalshi-mrb/backend/.venv/bin/uvicorn src.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
